# cron/.crontab.erb
SHELL=<%= ENV['SHELL'] %>
PATH=<%= ENV['PATH'] %>

@weekly zsh -c "source $HOME/.zshrc && backup-netvibes"

